# roles/mesos/master/files/etc/systemd/system/mesos-master.service
[Unit]
Description=Mesos Master
After=network.target

[Service]
ExecStart=/usr/bin/mesos-init-wrapper master

[Install]
WantedBy=multi-user.target
